# conf/supervisord.ini
[program:nanoinvoice]
command=/var/www/nanoinvoice.com/venv/bin/uwsgi --ini /var/www/nanoinvoice.com/conf/uwsgi.ini
directory=/var/www/nanoinvoice.com
autostart=true
autorestart=true
stdout_logfile=/var/www/nanoinvoice.com/log/uwsgi.log
redirect_stderr=true
stopsignal=QUIT
user=www-data
